# power_install/boottime_power.service
[Unit]
Description=Set power and wireless blocking at boot time

[Service]
ExecStart=/usr/bin/power boot

[Install]
WantedBy=multi-user.target 

